# infra/systemd/mc-bot.service
[Unit]
Description=Mission Control Telegram Bot
After=network.target mc-mcp.service
Requires=mc-mcp.service

[Service]
Type=simple
WorkingDirectory=%h/mission-control
ExecStart=%h/mission-control/.venv/bin/python -m agents.telegram_bot
Environment=PYTHONPATH=%h/mission-control
EnvironmentFile=%h/mission-control/.env
Restart=always
RestartSec=10
StandardOutput=append:%h/mission-control/logs/telegram-bot.log
StandardError=append:%h/mission-control/logs/telegram-bot.log

[Install]
WantedBy=default.target
