# infra/systemd/mc-bot.service
[Unit]
Description=Mission Control Telegram Bot
After=network.target mc-mcp.service
Requires=mc-mcp.service

[Service]
Type=simple
WorkingDirectory=/home/amaresh/projects/mission-control
ExecStart=/home/amaresh/projects/mission-control/.venv/bin/python -m mission_control.telegram_bot
Environment=PYTHONPATH=/home/amaresh/projects/mission-control
EnvironmentFile=/home/amaresh/projects/mission-control/.env
Restart=always
RestartSec=10
StandardOutput=append:/home/amaresh/projects/mission-control/logs/telegram-bot.log
StandardError=append:/home/amaresh/projects/mission-control/logs/telegram-bot.log

[Install]
WantedBy=default.target
