# infra/systemd/mc-scheduler.service
[Unit]
Description=Mission Control Scheduler
After=network.target mc-api.service

[Service]
Type=simple
WorkingDirectory={MC_HOME}
EnvironmentFile={MC_HOME}/.env
ExecStart={PYTHON} -m mission_control.scheduler_main
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
